# service/humer-sensors-[location].timer
[Unit]
Description=Humer timer - sensors readings - [location]

[Timer]
OnBootSec=5min
OnUnitActiveSec=[interval]

[Install]
WantedBy=timers.target