# service/humer-sensors-[location].timer
[Unit]
Description=Humer timer - sensors readings - [location]

[Timer]
OnCalendar=[interval]

[Install]
WantedBy=timers.target